# gmod/99-pwmdma.rules
ACTION=="add", KERNEL=="pwmdma", GROUP="dialout", MODE="660", SYMLINK+="dcc", RUN+="/bin/chgrp dialout /sys/devices/virtual/dccc/pwmdma/signal"
